# rucio/rucio-system.service
[Unit]
Description=Decentralized peer-to-peer file sharing over libp2p, with eMule/Kad2 compatibility
Documentation=https://github.com/ogarcia/rucio/blob/master/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=rucio
Environment=HOME=/var/lib/rucio
WorkingDirectory=/var/lib/rucio
ExecStart=/usr/bin/ruciod

[Install]
WantedBy=multi-user.target
