# polaris/polaris.service
[Unit]
Description=A self-hosted music streaming server
Documentation=https://github.com/agersant/polaris/blob/master/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=polaris
RuntimeDirectory=polaris
CacheDirectory=polaris
LogsDirectory=polaris
ExecStart=/usr/bin/polaris -f

[Install]
WantedBy=multi-user.target
