# rustic-server/rustic-server.service
[Unit]
Description=A REST server built in rust to use with rustic and restic
Documentation=https://rustic.cli.rs/ecosystem/rustic-server/
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=rustic-server
Group=rustic-server
ExecStart=/usr/bin/rustic-server -c /etc/rustic-server/rustic-server.toml serve
Restart=always

[Install]
WantedBy=multi-user.target
